# docker/visualize/supervisord.conf
[supervisord]
nodaemon=true

[program:visualize]
command=node /OSMDeepOD-Visualize/node-server/server.js
